# wifi-monitor_1.0.0-1/DEBIAN/postinst
#!/bin/sh

systemctl enable wifi-monitor.service
systemctl start wifi-monitor.service

exit 0

