# wifi-monitor_1.0.0-1/debian/postinst
#!/bin/sh

systemctl enable wifi-monitor.service
systemctl start wifi-monitor.service

systemctl enable wifi-led.service
systemctl start wifi-led.service

exit 0

