# wifi-monitor_1.0.0-1/DEBIAN/prerm
#!/bin/sh

systemctl stop wifi-monitor
systemctl disable wifi-monitor

systemctl stop wifi-led
systemctl disable wifi-led

exit 0

