# wifi-monitor_1.0.0-1/DEBIAN/prerm
#!/bin/sh

systemctl stop wifi-monitor
systemctl disable wifi-monitor

exit 0

